# deploy/epoch-server.service
[Unit]
Description=Epoch Snapshot Registry Server
After=network.target docker.service

[Service]
Type=simple
User=tiger
ExecStart=/home/tiger/bin/epoch serve --addr :4300 --dsn epoch:epoch@tcp(127.0.0.1:3306)/epoch?parseTime=true
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
